# supervisor/openresty.conf
[program:openresty]
command=/opt/openresty/nginx/sbin/nginx -p /var/www/ -c nginx.conf
autorestart=true
